# hub-pi/rfm69-gateway.service
[Unit]
Description=RFM69 to MQTT Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dietpi
WorkingDirectory=/home/dietpi/rfm69-pi-avr/hub-pi
ExecStart=/home/dietpi/rfm69-pi-avr/hub-pi/venv/bin/python hub.py --log-level=DEBUG
Restart=on-failure
RestartSec=10
ExecStartPre=/bin/sleep 5

[Install]
WantedBy=multi-user.target